# meta-auras/meta-cdu_g3-project/recipes-platform/common/modbus-tcp-server/src/service_files/xyz.openbmc_project.modbus_tcp_server.service
[Unit]
Description=Modbus TCP server
StopWhenUnneeded=false

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/modbus-tcp-server

[Install]
WantedBy=multi-user.target