# meta-auras/meta-cdu_g3-project/recipes-phosphor/sensors/stm32io/xyz.openbmc_project.ioService.service
[Unit]
Description=Auras STM32 IO service
StopWhenUnneeded=false
After=xyz.openbmc_project.EntityManager.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/ioService

[Install]
WantedBy=multi-user.target
